# packaging/systemd/user/shuvoice.service
[Unit]
Description=ShuVoice speech-to-text overlay
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
# Packaged install path. For repo/venv workflows, override with:
#   systemctl --user edit shuvoice.service
#   [Service]
#   ExecStart=
#   ExecStart=%h/.venv/bin/shuvoice
ExecStart=/usr/bin/shuvoice
Restart=on-failure
RestartSec=2
# Do not restart-loop when backend dependencies are missing.
RestartPreventExitStatus=78
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
